# scripts/infra/systemd/ibkr-bridge.service
[Unit]
Description=Agentyc IBKR Bridge (FastAPI + Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ibkr-bridge
ExecStart=/opt/ibkr-bridge/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target

